# crontab/run_seb_api_app.sh
#!/bin/sh
python3 /home/ricardas/test/seb_api_app/main.py
